# ec2/crontab
*/1 * * * * ~/Hey-Neighbor/script.sh start